# conf/celery.conf
start on started redis
stop on stopping redis

chdir /usr/local/fieldpapers/decoder
exec su ubuntu -c "screen -D -m -S celery -- celery -A tasks worker"

respawn
