# conf/celery.conf
start on (local-filesystems and net-device-up IFACE!=lo)

chdir /usr/local/fieldpapers/decoder
exec su ubuntu -c "screen -D -m -S celery -- celery -A tasks worker"

respawn
