# provisioning/configs/celery.conf
start on (local-filesystems and net-device-up IFACE!=lo)

env SENTRY_DSN=<...>

chdir /vagrant/fieldpapers/decoder
exec su ubuntu -c "screen -D -m -S celery -- celery -A tasks worker"

respawn
